# 40-gamepad.rules
ACTION=="add", SUBSYSTEMS=="input", RUN+="/bin/systemctl start gamepad.service"
